# playbooks/roles/automated_testing/templates/phantomjs.conf.j2
[program:phantomjs]
command=/usr/bin/phantomjs --webdriver=4445
autostart=true
autorestart=true
stderr_logfile=/var/log/phantomjs.err.log
stdout_logfile=/var/log/phantomjs.out.log