# playbooks/roles/automated_testing/templates/selenium.conf.j2
[program:selenium_start]
command=/usr/bin/selenium_start
autostart=true
autorestart=true
stderr_logfile=/var/log/selenium.err.log
stdout_logfile=/var/log/selenium.out.log