# playbooks/roles/automated_testing/templates/xvfb.conf.j2
[program:xvfb_start]
command=/usr/bin/xvfb_start
autostart=true
autorestart=true
stderr_logfile=/var/log/xvfb.err.log
stdout_logfile=/var/log/xvfb.out.log